# config/supervisor/app.conf
[program:app]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan queue:work redis --sleep=3 --tries=3
user=root
autostart=true
autorestart=true
log_stdout=true
log_stderr=true
